# ipa-enroll.service
[Unit]
Description=Enroll container to FreeIPA
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ipa-enroll.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
